# config/includes.chroot/etc/systemd/system/takesora-health.service
[Unit]
Description=TakesoraOS Health Score Reporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/takesora-health-score
