# config/includes.chroot/etc/systemd/system/takesora-cleaner.service
[Unit]
Description=TakesoraOS Auto Cleaner

[Service]
Type=oneshot
ExecStart=/usr/local/bin/takesora-cleaner
