# config/includes.chroot/etc/skel/.config/autostart/takesora-ai-panel.desktop
[Desktop Entry]
Type=Application
Name=Takesora AI Panel
Exec=gnome-terminal -- bash -lc 'echo "Takesora AI panel ready"; sleep 5'
X-GNOME-Autostart-enabled=true
